# deploy/systemd/relay-tunnel.service
[Unit]
Description=cloudflared dedicated tunnel for Relay (board/api/email.naples.agency)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared --no-autoupdate tunnel --config /root/.cloudflared/relay-config.yml run relay
Restart=always
RestartSec=5
User=root
WorkingDirectory=/root/.cloudflared
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
